# pcreporter.service
[Unit]
Description=pcreporter
ExecStart=pc_reporter
